# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:startup]
exitcodes=0
command=/bin/bash /opt/startup.sh

